# etc/systemd/system/prepare-docker.service
[Unit]
Before=docker.socket

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=-/usr/sbin/groupadd -r docker
